# roles/sentry/templates/sentry.crontab
0 3 * * * sentry SENTRY_CONF=/etc/sentry {{ sentry_root }}/env.d/bin/sentry cleanup --days=30
